# site-cookbooks/testapp/files/default/cc-server-worker.conf
start on started chef-solo
respawn
script
  cd /opt/cc-server/current
  ruby worker.rb
end script
