# site-cookbooks/testapp/files/default/cc-server.conf
start on started chef-solo
respawn
script
  cd /opt/cc-server/current
  bundle exec passenger start
end script
