# site-cookbooks/testapp/files/default/sidekiq.conf
start on started chef-solo
respawn

script
  cd /var/www/current
  bundle exec sidekiq -r ./app.rb -c 3
end script
